# pocket-cloud/systemd/pocketcloud-cleanup.timer
[Unit]
Description=PocketCloud Cleanup Timer
Documentation=https://github.com/yourname/pocketcloud
Requires=pocketcloud-cleanup.service

[Timer]
# Run daily at 3 AM
OnCalendar=*-*-* 03:00:00

# Run on boot if missed
Persistent=true

# Randomize start time by up to 30 minutes
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target